# windows-port-console/windows-port-console.service
[Unit]
Description=🖥️ Windows Port Console
Documentation=https://github.com/Baal-TehDriverman/vm-ai-gateway
After=network.target libvirtd.service
Wants=libvirtd.service

[Service]
Type=simple
User=tehlappy
Group=tehlappy
WorkingDirectory=/home/tehlappy/Desktop/windows-port-console
Environment=PATH=/home/tehlappy/Desktop/windows-port-console/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/tehlappy/Desktop/windows-port-console/venv/bin/python server.py
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGTERM
StandardOutput=append:/home/tehlappy/Desktop/windows-port-console/server.log
StandardError=append:/home/tehlappy/Desktop/windows-port-console/server.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/tehlappy/Desktop/windows-port-console /home/tehlappy/Desktop/cross-compile
ProtectHome=read-only

[Install]
WantedBy=multi-user.target